# scripts/metadata.cron
#! /bin/bash

/var/www/metadata/metadata backup