# scripts/mymetadata-go.cron
#! /bin/bash

/var/www/mymetadata-go/mymetadata-go backup